# packaging/s6/.s6-svscan/finish
#!/bin/sh

# No system-wide cleanup to run yet
