# packaging/s6/ws-test-01/run
#!/bin/sh

# If a no value has been set for the location of consul then run directly
if [ -z $CONSUL ]; then
    /usr/bin/ws-test-01
fi
